# tlb.upstart.in
# tlb - system startup/shutdown
#
# Copyright (c) 2023 Alessio Molinari <levnikmyskin at github.com> and others.
# This software is licensed under the GPL v2 or later.

description "tlb"

start on ( virtual-filesystems and runlevel [2345] )
stop on runlevel [!2345]

env TLB=@TLB_SBIN@/tlb

pre-start script
    [ -x $TLB ] || exit 4
    $TLB init start

end script


