# tlb.service.in
# tlb - systemd startup/shutdown service
#
# Copyright (c) 2023 Alessio Molinari <levnikmyskin at github.com> and others.
# SPDX-License-Identifier: GPL-2.0-or-later

[Unit]
Description=TLB system startup/shutdown
After=multi-user.target NetworkManager.service
Before=shutdown.target
Documentation=https://linrunner.de/tlb

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@TLB_SBIN@/tlb init start
ExecReload=@TLB_SBIN@/tlb start

[Install]
WantedBy=multi-user.target
